# server_config/services/biometric_ingest/biometric-router-worker.service
[Unit]
Description=Biometric Event Router Worker
After=network.target postgresql.service biometric-ingest.service
Wants=postgresql.service biometric-ingest.service

[Service]
Type=simple
User=biometric
Group=biometric
WorkingDirectory=/opt/biometric-ingest
EnvironmentFile=/etc/biometric-ingest.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/biometric-ingest/venv/bin/python /opt/biometric-ingest/biometric_router_worker.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/lib/biometric-ingest /var/log/biometric-ingest /var/log/dahua_events.log

[Install]
WantedBy=multi-user.target
